# deploy/scan-button/papermind-scan-idle.service
[Unit]
Description=PaperMind Scan – Batch automatisch abschließen, wenn er ruht
# Sicherheitsnetz, falls die "PDF/file"-Taste mal vergessen wird.

[Service]
Type=oneshot
Environment=IDLE_SECONDS=180
# Pfad an die tatsächliche Repo-Position anpassen (scan-inbox wird vom Script
# selbst aus dem Repo-Pfad abgeleitet):
ExecStart=/home/jan/papermind/deploy/scan-button/papermind-scan.sh finalize-idle
# Gleiche Identität wie scanbd (siehe scanbd.conf user/group), damit Lock/Batch
# und Dateibesitz konsistent sind:
User=jan
Group=jan
