# deploy/scan-button/papermind-scan-watch.service
[Unit]
Description=PaperMind Scan-Tasten-Poller (Canon LiDE 400, Ersatz für scanbd)
Wants=papermind-scanner-usb-awake.service
After=multi-user.target papermind-scanner-usb-awake.service

[Service]
Type=simple
# Pfad an die tatsächliche Repo-Position anpassen:
ExecStart=/home/jan/papermind/deploy/scan-button/papermind-scan-watch.sh
Restart=always
RestartSec=5
# Muss der Eigentümer von scan-inbox sein und den Scanner nutzen dürfen:
User=jan
Group=jan
# Optional festes Gerät statt Auto-Erkennung:
# Environment=SCAN_DEVICE=pixma:04A91912_50F25C
# Adaptives Poll-Intervall: haelt den USB-Scanner reaktionsbereit.
Environment=ACTIVE_POLL_INTERVAL=0.35
Environment=IDLE_POLL_INTERVAL=1
Environment=ACTIVE_WINDOW_SECONDS=30

[Install]
WantedBy=multi-user.target
